# openclaw_agents/operations/systemd/openclaw-worker-supervisor.service
[Unit]
Description=OpenClaw runtime worker supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/workspace/claw_software_workspace/.agents/state/openclaw_agents/env/openclaw-runtime-workers.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/env bash -lc 'cd "$$OPENCLAW_REPO_ROOT" && cmd=(python3 -m openclaw_agents.runtime.worker_supervisor --config "$$OPENCLAW_WORKER_CONFIG" --check); if [[ -n "$$OPENCLAW_RUNTIME_STATE_DIR" ]]; then cmd+=(--state-dir "$$OPENCLAW_RUNTIME_STATE_DIR"); fi; "$${cmd[@]}"'
ExecStart=/usr/bin/env bash -lc 'cd "$$OPENCLAW_REPO_ROOT" && cmd=(python3 -m openclaw_agents.runtime.worker_supervisor --config "$$OPENCLAW_WORKER_CONFIG"); if [[ -n "$$OPENCLAW_RUNTIME_STATE_DIR" ]]; then cmd+=(--state-dir "$$OPENCLAW_RUNTIME_STATE_DIR"); fi; exec "$${cmd[@]}"'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
