# openclaw_agents/systemd/zulip-gateway-v3.service
[Unit]
Description=Zulip Gateway V3
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=LOCAL_BRIDGE_USER
Group=LOCAL_BRIDGE_USER
WorkingDirectory=LOCAL_ZULIP_GATEWAY_V3_DIR
ExecStart=LOCAL_ZULIP_GATEWAY_V3_DIR/run_gateway.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
